# Docker/logstash/supervisord.conf
[supervisord]
nodaemon = true

[program:logstash]
command = /usr/share/logstash/bin/logstash
